# loop/loop-alloc.service
[Unit]
Description=Loop device allocator

[Service]
ExecStart=/usr/local/bin/loop-alloc

[Install]
WantedBy=libvirtd.service
